# package/ubnt-base-files/files/usr/etc/hotplug.d/usb/10-usb
#!/bin/sh
case "$ACTION" in
	add)
		;;
	remove)
		;;
esac    
